# Linux/MySQL/dump-database.timer
[Unit]
Description=Run every ten minutes

[Timer]
OnCalendar=*-*-* *:00,10,20,30,40,50:00

[Install]
WantedBy=timers.target
