# Linux/MySQL/dump-database.service
[Unit]
Description=Dump mysql database

[Service]
ExecStart=/usr/bin/dump-database.sh

[Install]
WantedBy=multi-userr.target